# cron/crontab
*/30 * * * * /bin/sh app/run_scan.sh >> /var/log/cron.log 2>&1

