# cron/crontab
*/30 * * * * /bin/sh /scripts/run_scan.sh >> /var/log/cron.log 2>&1